# crontab
15 * * * * root python /update_cloudflare.py -u
